# ubuntu/crontab/debitoautomatico.sh
#!/bin/bash
set -e
cd ../..
. venv/bin/activate
exec python manage.py notificar
